# bootscripts/sensorReader
### BEGIN INIT INFO
# Provides: LCD - data
# Required-Start: $remote_fs $syslog
# Required-Stop: $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Sensor loggin and alarm printer
# Description: data-push / LCD / piezo / LED
### END INIT INFO
 
 
#! /bin/sh
# /etc/init.d/sensorReader
 
 
export HOME
case "$1" in
    start)
        echo "Starting monitoring"
        /home/pi/UNICEF_VACCINE_MONITOR/sensorReader.py 2>&1 &
    ;;
    stop)
        echo "Stopping monitoring"
		reader_PID=`ps auxwww | grep sensorReader.py | head -1 | awk '{print $2}'`
		kill -9 $reader_PID
    ;;
    *)
        echo "Usage: /etc/init.d/sensorReader {start|stop}"
        exit 1
    ;;
esac
exit 0
